# misc/desktop_autostart/advvcom_dkms.desktop
[Desktop Entry]
Type=Application
Name=Advantech VCOM DKMS status check
Exec=/usr/local/advtty/adv_dkms_stat_check
X-GNOME-Autostart-enabled=true
